# scripts/trimix-analyzer.service
[Unit]
Description=Trimix Analyzer - Native RPi Gas Analyzer
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/trimix-analyzer
Environment="PATH=/opt/trimix-analyzer/.venv/bin"
Environment="PYTHONPATH=/opt/trimix-analyzer"
Environment="TRIMIX_ENVIRONMENT=production"
Environment="TRIMIX_MOCK_SENSORS=0"
ExecStart=/opt/trimix-analyzer/.venv/bin/python main.py
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trimix-analyzer

# Security and resource limits
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
MemoryMax=400M

# GPIO and I2C access
SupplementaryGroups=gpio i2c spi video

[Install]
WantedBy=multi-user.target
